# src/daily_flight_deck/scheduler_templates/com.rory.daily-flight-deck.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.rory.daily-flight-deck</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/roryh/Documents/Rorys Dashboard/scripts/run_snapshot_job.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/roryh/Documents/Rorys Dashboard</string>

  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>6</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>30</integer></dict>
  </array>

  <key>StandardOutPath</key>
  <string>/Users/roryh/Documents/Rorys Dashboard/data/scheduler.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/roryh/Documents/Rorys Dashboard/data/scheduler.err.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
